# cron-caller/crontab
*/20 * * * * /usr/local/bin/call_api.sh >> /var/log/cron.log 2>&1
